# extras/com.mt-butterfly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mt-butterfly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/.local/bin/mt-butterfly</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/mt-butterfly/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/mt-butterfly/stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME</string>
</dict>
</plist>
